# config/launchd/com.ai-company-os.daily-operating-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-company-os.daily-operating-loop</string>

    <!-- Schedule: daily at 09:00 (system timezone = CST UTC+8) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Run once per trigger — don't linger if previous run hung -->
    <key>ThrottleInterval</key>
    <integer>3600</integer>

    <!-- Working directory = project root, so relative paths in scripts work -->
    <key>WorkingDirectory</key>
    <string>/Users/tangbomao/Documents/Codex/ai-company-os</string>

    <!-- Use project's Python from Hermes venv — launchd PATH is minimal -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/tangbomao/.hermes/hermes-agent/venv/bin/python3</string>
        <string>scripts/run_operating_loop.py</string>
        <string>--once</string>
        <string>--wait-results</string>
        <string>--timeout</string>
        <string>120</string>
    </array>

    <!-- Environment — prevent dual-DB bugs -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>AI_COMPANY_OS_DATABASE_PATH</key>
        <string>/Users/tangbomao/Documents/Codex/ai-company-os/backend/data/ai_company_os.db</string>
        <key>OPENCLAW_EXECUTOR_MODE</key>
        <string>backend</string>
        <key>PATH</key>
        <string>/Users/tangbomao/.hermes/hermes-agent/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Logs — stdout and stderr separated -->
    <key>StandardOutPath</key>
    <string>/Users/tangbomao/Documents/Codex/ai-company-os/logs/launchd/daily-operating-loop.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tangbomao/Documents/Codex/ai-company-os/logs/launchd/daily-operating-loop.err.log</string>

    <!-- Keep trying on transient failures (network, DB locked) -->
    <key>KeepAlive</key>
    <false/>

    <!-- Don't run at load time, only on schedule -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
