# scripts/com.openclaw.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.heartbeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/tangbomao/.openclaw/workspace/heartbeat.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/tangbomao/.openclaw/workspace/memory/heartbeat.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tangbomao/.openclaw/workspace/memory/heartbeat.launchd.error.log</string>
</dict>
</plist>
